# systemd/corridor-data.service.in
[Unit]
Description=corridor's relay list
Documentation=man:corridor(8)
ConditionPathExists=|/var/run/qubes-service/corridor
ConditionPathExists=|!/var/run/qubes-service/

[Service]
ExecStart=@SBIN@/corridor-data
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
