# systemd/corridor-init-snat.service.in
[Unit]
Description=corridor's source NAT
Documentation=man:corridor(8)
After=iptables.service
ConditionPathExists=|/var/run/qubes-service/corridor
ConditionPathExists=|!/var/run/qubes-service/

[Service]
ExecStart=@SBIN@/corridor-init-snat
ExecStop=@SBIN@/corridor-stop-snat
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
